# infra/techblogs-scraper.service
[Unit]
Description=Techblogs scraper (one-shot)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=techblogs
Group=www-data
Environment=DB_PATH=/srv/techblogs/data/techblogs.db
Environment=SQLITE_TMPDIR=/srv/techblogs/data
ExecStart=/srv/techblogs/current/backend/techblogs-scraper
WorkingDirectory=/srv/techblogs

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/srv/techblogs/data
PrivateTmp=true
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallFilter=@network-io @system-service
CapabilityBoundingSet=
SystemCallFilter=@system-service
MemoryDenyWriteExecute=true
LockPersonality=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
CPUQuota=50%
MemoryMax=500M
TimeoutSec=300

[Install]
WantedBy=multi-user.target
